# repair_ro
#!/bin/sh /etc/rc.common
# Created by Frizkyiman
# Github https://github.com/frizkyiman/

START=99

start() {
    if [ ! -d /usr/bin/repair_ro ]; then
        cat <<'EOF' >/usr/bin/repair_ro
#!/bin/sh
/etc/init.d/repair_ro start
EOF
        chmod +x /usr/bin/repair_ro
    fi

    echo "Repairing file system..."

    # Check if the file system is currently in read-only mode
    if grep -qs " ro " /proc/mounts; then
        # Find the root file system device
        root_device=$(cat /proc/cmdline | grep -o 'root=[^ ]*' | sed 's/root=//')

        if [ -z "$root_device" ]; then
            echo "Root file system not found. Unable to repair."
            return
        fi

        echo "Detected root file system: $root_device"

        echo "The file system is currently read-only. Initiating repair..."

        # Perform file system repair using e2fsck -y
        e2fsck -y "$root_device"

        # Check if the repair was successful
        if [ $? -eq 0 ]; then
            echo "File system repair successful."
            echo "Rebooting the system..."
            reboot
        else
            echo "An error occurred while repairing the file system."
        fi

        # Restore the file system to read-write mode
        mount -o remount,rw "$root_device" /
    else
        echo "The file system is currently in read-write mode."
    fi
}
