# etc/init.d/repair_ro
#!/bin/sh /etc/rc.common
# Created by Frizkyiman
# Github https://github.com/frizkyiman/

START=01

start() {
    # Finding root partition
    root_device=$(cat /proc/cmdline | grep -o 'root=[^ ]*' | sed 's/root=//')

    if [ -z "$root_device" ]; then
        echo "Root file system not found. Unable to repair."
        exit 1
    fi

    # Checking root partition status
    if mount | grep -q " / " | grep -q "ro"; then
        echo "Detected root file system: $root_device"
        echo "The file system is currently read-only. Initiating repair..."

        e2fsck -y "$root_device"

        if [ $? -eq 0 ]; then
            echo "File system repair successful."
            echo "Rebooting the system..."
            reboot
        else
            echo "An error occurred while repairing the file system."
        fi

        mount -o remount,rw "$root_device" / 
    else
        echo "The root file system "$root_device" is currently in read-write mode."
    fi
    }
